# backend/supervisord.conf
; The supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Quotes around values are not supported, except in the case of
;    the environment= options as shown below.
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".
;  - Command will be truncated if it looks like a config file comment, e.g.
;    "command=bash -c 'foo ; bar'" will truncate to "command=bash -c 'foo ".

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/logs/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
;pidfile=pids/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:ligiermirror]
command=ligiermirror -m %(ENV_TAGS_TO_MIRROR)s -t monitoring_ligier_1 -p %(ENV_DAQ_LIGIER_PORT)s %(ENV_DAQ_LIGIER_IP)s
priority=10
stdout_logfile=/logs/%(program_name)s.out.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/logs/%(program_name)s.err.log        ; stderr log path, NONE for none; default AUTO

[program:acoustics]
command=python -u scripts/acoustics.py -d %(ENV_DETECTOR_ID)s
stdout_logfile=/logs/%(program_name)s.out.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/logs/%(program_name)s.err.log        ; stderr log path, NONE for none; default AUTO

[program:ahrs_calibration]
command=python -u scripts/ahrs_calibration.py -d %(ENV_DETECTOR_ID)s -l monitoring_ligier_1
;process_name=%(program_name)s ; process_name expr (default %(program_name)s)
;priority=999                  ; the relative start priority (default 999)
;startsecs=1                   ; # of secs prog must stay up to be running (def. 1)
;startretries=3                ; max # of serial start failures when starting (default 3)
;autorestart=unexpected        ; when to restart if exited after running (def: unexpected)
;redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/logs/%(program_name)s.out.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/logs/%(program_name)s.err.log        ; stderr log path, NONE for none; default AUTO

[program:dom_activity]
command=python -u scripts/dom_activity.py -d %(ENV_DETECTOR_ID)s -l monitoring_ligier_1
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

[program:dom_rates]
command=python -u scripts/dom_rates.py -d %(ENV_DETECTOR_ID)s -l monitoring_ligier_1
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

[program:pmt_rates]
command=python -u scripts/pmt_rates.py -d %(ENV_DETECTOR_ID)s -l monitoring_ligier_1 -i 20 -u 3
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

[program:trigger_rates]
command=python -u scripts/trigger_rates.py -l monitoring_ligier_1
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log
 
[program:triggermap]
command=python -u scripts/live_triggermap.py -d %(ENV_DETECTOR_ID)s -l monitoring_ligier_1
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

[program:ztplot]
command=python -u scripts/ztplot.py -d %(ENV_DETECTOR_ID)s -l monitoring_ligier_1
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

;[program:royfit]
;command=julia scripts/live_royfit.jl %(ENV_DETX)s %(ENV_MONITORING_LIGIER_PORT)s %(ENV_ROYFIT_TIMERES)s
;stdout_logfile=/logs/%(program_name)s.out.log
;stderr_logfile=/logs/%(program_name)s.err.log

;[program:time_residuals]
;command=python -u scripts/time_residuals.py %(ENV_ROYFIT_TIMERES)s
;stdout_logfile=/logs/%(program_name)s.out.log
;stderr_logfile=/logs/%(program_name)s.err.log

[program:rttc]
command=python -u scripts/rttc.py -d %(ENV_DETECTOR_ID)s -l %(ENV_DETECTOR_MANAGER_IP)s
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

[program:acoustics]
command=python -u scripts/acoustics.py -d %(ENV_DETECTOR_ID)s
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

;[program:k40]
;command=python -u scripts/k40_calibration.py -d %(ENV_DETECTOR_ID)s -p %(ENV_MONITORING_LIGIER_PORT)s
;stdout_logfile=/logs/%(program_name)s.out.log
;stderr_logfile=/logs/%(program_name)s.err.log

[program:timesync_monitor]
command=python -u scripts/timesync_monitor.py -l monitoring_ligier_1 -m %(ENV_LOG_LIGIER_IP)s -q %(ENV_LOG_LIGIER_PORT)s
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

[program:msg_dumper]
command=python -u scripts/msg_dumper.py -l %(ENV_LOG_LIGIER_IP)s -p %(ENV_LOG_LIGIER_PORT)s -o /logs
priority=200
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

[program:log_analyser]
command=python -u scripts/log_analyser.py
priority=205
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

[program:chatbot]
command=python -u scripts/chatbot.py
priority=220
stdout_logfile=/logs/%(program_name)s.out.log
stderr_logfile=/logs/%(program_name)s.err.log

[group:ligiers]
programs=ligiermirror
priority=1

[group:logging]
programs=msg_dumper,log_analyser,chatbot
priority=200

[group:monitoring_process]
programs=acoustics,ahrs_calibration,dom_activity,dom_rates,pmt_rates,trigger_rates,triggermap,ztplot,rttc
priority=500

[group:alerts]
programs=timesync_monitor
priority=600

;[group:reconstruction]
;programs=royfit,time_residuals
;priority=999
